# profile.d/x11.sh
#!/bin/sh

export TERMINAL='st'
export DMENU='dmenu'

export SXHKD_SHELL=/bin/sh
