# profile.d/x11.sh
#!/bin/sh

export TERMINAL=$(command -v st)
export DMENU=$(command -v dmenu)

export SXHKD_SHELL=/bin/sh
